# configs/crontab
*/5 * * * * php -f /var/www/sync.php >/dev/null 2>&1
